# carousell-uploader.service
[Unit]
Description=Carousell Uploader Service
After=network.target

[Service]
Type=simple
User=carousell
Group=carousell
WorkingDirectory=/opt/carousell_upload
Environment=PATH=/opt/carousell_upload/venv/bin
ExecStart=/opt/carousell_upload/venv/bin/python -m cli.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/carousell-uploader/logs /opt/carousell-uploader/data

[Install]
WantedBy=multi-user.target
